# enrollment.service
[Unit]
Description=Enrollment service with docker-compose
Requires=docker.service
After=docker.service

[Service]
Restart=always

WorkingDirectory=/usr/src/enrollment

# Remove old containers, network and volumes
ExecStartPre=/usr/local/bin/docker-compose -f docker/docker-compose.yaml down -v
ExecStartPre=-/bin/bash -c 'docker network rm docker_jfrog_net'
ExecStartPre=-/bin/bash -c 'docker ps -aqf "name=artifactory" | xargs docker rm'

# Compose up
ExecStart=/usr/local/bin/docker-compose -f docker-compose.yaml up -d

# Compose down, remove containers
ExecStop=/usr/local/bin/docker-compose -f docker-compose.yaml down

[Install]
WantedBy=multi-user.target